# ansible-pull.service
[Unit]
Description=Ansible pull

[Service]
WorkingDirectory=/var/lib/ansible/pull
ExecStart=/usr/bin/ansible-playbook playbook.yml
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
